# launchd/ai.radxa.mascot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.radxa.mascot</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ProgramArguments</key>
    <array>
      <string>__HOME__/.radxa-mascot/mac_mini/run_openclaw_server.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__HOME__/.radxa-mascot</string>
    <key>EnvironmentVariables</key>
    <dict>
      <key>MASCOT_TTS</key>
      <string>openclaw_tts</string>
      <key>MASCOT_SAY_VOICE</key>
      <string>Kyoko</string>
      <key>MASCOT_TTS_VOICE</key>
      <string>ja-JP-NanamiNeural</string>
      <key>MASCOT_TTS_TIMEOUT</key>
      <string>90</string>
      <key>MASCOT_PUBLIC_BASE_URL</key>
      <string>__MASCOT_PUBLIC_BASE_URL__</string>
      <key>OPENCLAW_MASCOT_MODEL</key>
      <string>openai/gpt-5.4-mini</string>
      <key>OPENCLAW_MASCOT_BACKEND</key>
      <string>infer</string>
      <key>OPENCLAW_MASCOT_TIMEOUT</key>
      <string>45</string>
      <key>MASCOT_STT_MODEL</key>
      <string>mlx-community/whisper-tiny-mlx</string>
      <key>MASCOT_STT_DIRECT</key>
      <string>on</string>
      <key>MASCOT_STT_TIMEOUT</key>
      <string>30</string>
      <key>MASCOT_HOST</key>
      <string>0.0.0.0</string>
      <key>MASCOT_PORT</key>
      <string>8765</string>
      <key>PATH</key>
      <string>/opt/homebrew/bin:__HOME__/.npm-global/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/radxa-mascot/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/radxa-mascot/stderr.log</string>
  </dict>
</plist>
